# profile.d/01-path
#!/bin/bash

export PATH

path_push() {
  PATH="${PATH/":$1"}:$1"
}
path_unshift() {
  PATH="$1:${PATH/"$1:"}"
}

[ -d "$HOME/bin" ] && path_push "$HOME/bin"
