# profile.d/30-linux
#!/bin/bash

path_unshift "/usr/local/bin"
alias ls='ls --color=auto'
