# profile.d/30-mac
#!/bin/bash

# Mac-specific things
alias lock='/System/Library/CoreServices/Menu\ Extras/User.menu/Contents/Resources/CGSession -suspend'
