# profile.d/40-golang
#!/bin/bash
export GOPATH="$HOME/go"
export PATH="$PATH:$GOPATH/bin"

mkjump gogogo "$HOME/go/src/github.com/zenazn" <<EOF
EOF
