# profile.d/40-golang
#!/bin/bash
export GOPATH="$HOME/go"
add_to_path "$GOPATH/bin"

mkjump z "$HOME/go/src/github.com/zenazn" <<EOF
EOF
