# profile.d/40-golang
#!/bin/bash
export GOPATH="$HOME/go"
path_push "$GOPATH/bin"

mkjump z "$HOME/go/src/github.com/zenazn" <<EOF
EOF

alias gotags="gotags -R . > tags"
